# lab2-bonus/crontab
* * * * * sh  /home/asala/6916-lab2/lab2-bonus/backup-cron.sh /home/asala/6916-lab2/lab2-bonus/dir /home/asala/6916-lab2/lab2-bonus/backup 2
